# installer/dataease/bin/dataease/dataease.service
[Unit]
Description=DataEase Service
After=docker.service
Requires=docker.service

[Service]
User=root
Group=root

Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/dectl start
ExecStop=/usr/bin/dectl stop
ExecReload=/usr/bin/dectl reload

[Install]
WantedBy=multi-user.target